# chef/site-cookbooks/mycookbook/templates/default/play.erb
#!/bin/bash
# chkconfig: 345 20 80
# description: Play start/shutdown script
# processname: play
#
# Instalation:
# copy file to /etc/init.d
# chmod +x /etc/init.d/play
# chkconfig --add /etc/init.d/play
# chkconfig play on
#
# Usage: (as root)
# service play start
# service play stop
# service play status
#
# Remember, you need python 2.6 to run the play command, it doesn't come standard with RedHat/Centos 5.5
# Also, you may want to temporarily remove the >/dev/null for debugging purposes
 
# Path to play install folder
PLAY_HOME=<%= @play_home %>
PLAY=$PLAY_HOME/play
 
# Path to the JVM
JAVA_HOME=/usr/lib/jvm/java
export JAVA_HOME
 
# User running the Play process
USER=<%= @user %>
 
# Path to the application
APPLICATION_PATH=<%= @application_path %>
APPLICATION_MODE=prod
 
#APPLICATION2_PATH=/path/to/application2
#APPLICATION_MODE=prod
 
# source function library
. /etc/init.d/functions
RETVAL=0
 
start() {
	if [ -e ${APPLICATION_PATH}/server.pid ]; then
		echo -n "Already running "
		echo_failure
		echo
		exit 1
	fi

	echo -n "Starting Play service: "
	su -s /bin/sh $USER -c "${PLAY} start ${APPLICATION_PATH} --%${APPLICATION_MODE} >/dev/null"
	RETVAL=$?
	
	# You may want to start more applications as follows
	# [ $RETVAL -eq 0 ] && su -s /bin/sh $USER -c "${PLAY} start ${APPLICATION2_PATH} --%${APPLICATION_MODE} > /dev/null"
        # RETVAL=$?
 
	if [ $RETVAL -eq 0 ]; then
		echo_success
	else
		echo_failure
	fi
	echo
}
stop() {
	echo -n "Shutting down Play service: "
	${PLAY} stop ${APPLICATION_PATH} > /dev/null
	# ${PLAY} stop ${APPLICATION2_PATH} > /dev/null
	
	RETVAL=$?
	
	if [ $RETVAL -eq 0 ]; then
		echo_success
	else
		echo_failure
	fi
	echo
}
status() {
	${PLAY} status ${APPLICATION_PATH}
	#${PLAY} status ${APPLICATION2_PATH}
	RETVAL=$?
 
	if [ $RETVAL -eq 0 ]; then
		echo "play is running"
		exit 0
	else
		echo "play is stopped"
		exit 1
	fi
}
clean() {
        rm -f ${APPLICATION_PATH}/server.pid
        #rm -f ${APPLICATION2_PATH}/service.pid
}
case "$1" in
	start)
	start
	;;
	stop)
	stop
	;;
	restart|reload)
	stop
	sleep 3
	start
	;;
	status)
	status
	;;
	clean)
	clean
	;;
	*)
	echo "Usage: $0 {start|stop|restart|status}"
esac
exit 0
